# rootfs/etc/services.d/aria2-ban/finish
#!/usr/bin/with-contenv bash

bash /Aria2-Ban/shutdown.sh
