# rootfs/etc/services.d/aria2-ban/run
#!/usr/bin/with-contenv bash

LOCAL_IP=$(cat /etc/hostname)
aria2b -u "http://$LOCAL_IP:$RPC_PORT/jsonrpc" -s "$RPC_SECRET"